# competitor-monitor/com.vantagecircle.competitor-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.vantagecircle.competitor-monitor</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/sahil/Library/Application Support/VantageCompetitorMonitor/run_monitor.py</string>
  </array>

  <!-- Run every Monday at 11:00 -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>1</integer>
    <key>Hour</key>
    <integer>11</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/sahil/Library/Application Support/VantageCompetitorMonitor/logs/launchd-stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/sahil/Library/Application Support/VantageCompetitorMonitor/logs/launchd-stderr.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
